# service/systemd/mosquitto.socket
[Socket]
ListenStream=/run/mosquitto/s
ListenStream=127.0.0.1:1883
ListenStream=[::1]:1883

[Install]
WantedBy=sockets.target
